# systemd/niralisd.service
[Unit]
Description=Niralis Display Manager Daemon
Documentation=man:niralisd(8)
After=systemd-user-sessions.service

[Service]
Type=simple
# Keep the inheritable capability set empty before Niralis starts. The worker
# deliberately rejects any capability inherited across its privileged child
# boundary; permitted/effective capabilities remain available until that child
# performs the verified credential drop.
ExecStart=/usr/bin/setpriv --inh-caps=-all -- /usr/sbin/niralisd --config /etc/niralis/niralis.toml
Restart=on-failure
RuntimeDirectory=niralis
RuntimeDirectoryMode=0700
UMask=0077

[Install]
WantedBy=graphical.target
